# llama-server.service
[Unit]
Description=llama.cpp Server Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/systemd/user/llama-server.env
ExecStart=/bin/bash -c 'source %h/.config/systemd/user/llama-server.env; if [ -f "${MODEL_PATH}" ]; then MODEL_ARG="-m"; else MODEL_ARG="-hf"; fi; exec "${LLAMA_SERVER_BIN}" "${MODEL_ARG}" "${MODEL_PATH}" --host "${HOST}" --port "${PORT}" $([ -n "${NGL_LEVEL}" ] && echo "-ngl \"${NGL_LEVEL}\"" || true) --ctx-size "${CONTEXT_SIZE}" $([ "${JINJA_ENABLED}" = "true" ] && echo "--jinja" || true)'
Restart=always
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
